# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[program:pioinit]
command=/pioinit.sh
autostart=true
autorestart=false

[program:elasticsearch]
command=/root/PredictionIO/vendors/elasticsearch-1.4.2/bin/elasticsearch
priority=1
startsecs=0
stdout_events_enabled=true
stderr_events_enabled=true
autostart=false
autorestart=false

[program:hbase]
command=/root/PredictionIO/vendors/hbase-0.98.6/bin/hbase --config /root/PredictionIO/vendors/hbase-0.98.6/conf master start
priority=2
startsecs=0
stdout_events_enabled=true
stderr_events_enabled=true
autorestart=true
autostart=false

[program:eventserver]
command=/root/PredictionIO/bin/pio eventserver --ip 0.0.0.0
priority=3
stdout_events_enabled=true
stderr_events_enabled=true
startsects=0
autoresart=true
autostart=false

[eventlistener:stdout]
command = supervisor_stdout
buffer_size = 100
events = PROCESS_LOG
result_handler = supervisor_stdout:event_handler
